# deploy/systemd/aigs.service
[Unit]
Description=AIGS - AI Gateway Service
Documentation=https://github.com/your-org/aigs
After=network.target
Wants=network.target

[Service]
Type=notify
NotifyAccess=all
ExecStart=/opt/aigs/Aigs.Api
WorkingDirectory=/opt/aigs
User=aigs
Group=aigs

# Environment variables
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000
EnvironmentFile=/etc/aigs.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aigs/data /opt/aigs/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictNamespaces=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aigs

[Install]
WantedBy=multi-user.target
